# deploy/news-aggregator.service
[Unit]
Description=Personal AI Crypto Research Aggregator
After=network-online.target
Wants=network-online.target

[Service]
# Type=notify lets us tell systemd READY=1 once the bot + scheduler are up,
# and ping WATCHDOG=1 from the event loop. WatchdogSec converts the
# "process alive but event loop wedged" failure mode (2026-05-28 incident)
# into an automatic restart, which Restart=on-failure cannot do on its own
# because a wedged process never exits.
Type=notify
NotifyAccess=main
WatchdogSec=180
User=news-bot
Group=news-bot
WorkingDirectory=/opt/news-aggregator
EnvironmentFile=/opt/news-aggregator/.env
Environment=NEWS_AGGREGATOR_DATA_DIR=/var/lib/news-aggregator
ExecStart=/opt/news-aggregator/.venv/bin/python -m aggregator
Restart=on-failure
RestartSec=10s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=news-aggregator

# Hardening. Zero-cost defense in depth; the bot only needs egress HTTPS
# and write access to its SQLite data dir.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/news-aggregator
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target
